# cistern-arch-designer.service
[Unit]
Description=Cistern Arch Designer Web UI
After=network.target

[Service]
Type=simple
ExecStart=/home/lobsterdog/go/bin/arch-designer --web --port 5738
Restart=on-failure
RestartSec=5
User=lobsterdog
StandardOutput=append:/home/lobsterdog/.cistern/arch-designer.log
StandardError=append:/home/lobsterdog/.cistern/arch-designer.log

[Install]
WantedBy=multi-user.target
